# cron/crontab
# === user crontab (managed by dotfiles) ===

# Rotate X session logs
0 3 * * * /usr/sbin/logrotate $HOME/.config/logrotate/xsession \
  >> $HOME/.local/state/xsession/logrotate.log 2>&1
